# systemd/funlog.service
[Unit]
Description=Run my fun data logging timer

[Service]
Type=oneshot
ExecStart=/home/calum/.local/custBin/log-fun-data
